# chromeos-base/tty/files/tty-base.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Terminal login service for VT"
author          "chromium-os-dev@chromium.org"

start on started boot-services
stop on stopping boot-services
respawn

env TTY_BAUD_RATE=1500000

# This file is generated at build time by chromeos-base/tty.
script
  if crossystem "cros_debug?1"; then
    exec agetty "${TTY_BAUD_RATE}" %PORT% linux
  else
    stop
    exit 0
  fi
end script
